# agents/linux/guardian-agent.service
[Unit]
Description=Guardian Agent Enterprise
After=network.target

[Service]
Type=simple
User=guardian-agent
WorkingDirectory=/opt/guardian-agent
ExecStart=/usr/bin/python3 /opt/guardian-agent/guardian_agent.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
